# pi/systemd/robo-serial-ingestor.service
# Copie para /etc/systemd/system/robo-serial-ingestor.service e ajuste
# User, WorkingDirectory e o caminho do .venv conforme onde o repo está.
#
#   sudo cp pi/systemd/robo-serial-ingestor.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now robo-serial-ingestor
#   journalctl -u robo-serial-ingestor -f   # ver logs

[Unit]
Description=Robo PIE V - Serial Ingestor (ESP32 -> MQTT)
# Sobe depois da rede e do broker local em container.
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
# Troque para o usuário do projeto (precisa estar no grupo 'dialout').
User=setrem
WorkingDirectory=/home/setrem/setrem-robot/pi/services/serial_ingestor
ExecStart=/home/setrem/setrem-robot/pi/services/serial_ingestor/.venv/bin/serial-ingestor
Environment=SERIAL_PORT=/dev/ttyUSB0
Environment=SERIAL_BAUD=115200
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883
# Reinicia sempre que cair; espera 3s entre tentativas.
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
