# pi/systemd/robo-wifi.service
# Copie para /etc/systemd/system/robo-wifi.service e ajuste os caminhos
# conforme onde o repo está.
#
#   sudo cp pi/systemd/robo-wifi.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now robo-wifi
#   journalctl -u robo-wifi -f   # ver logs
#
# Roda como ROOT de propósito: precisa aplicar a configuração de rede com o
# nmcli. NÃO faz Bluetooth — a credencial chega pelo barramento MQTT (o app
# manda via ESP32 -> serial). Por isso depende do NetworkManager.

[Unit]
Description=Robo PIE V - Wi-Fi provisioning (comando MQTT -> nmcli)
After=network-online.target docker.service NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/setrem/setrem-robot/pi
ExecStart=/home/setrem/setrem-robot/pi/.venv/bin/wifi
Environment=WIFI_IFACE=wlan0
Environment=MQTT_HOST=127.0.0.1
Environment=MQTT_PORT=1883
# Reinicia sempre que cair; espera 3s entre tentativas.
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
